# files/etc/periodic/hourly/get_iplayer_pvr
#!/bin/sh
#
# Run pvr
#
umask 2
su-exec $PUID:$PGID /usr/local/bin/get_iplayer --pvr --quiet 2>> /var/log/get_iplayer.log
